# supervisor.c3-pro-idm-frontend.conf
[program:c3-pro-idm-frontend]
directory   = /home/ubuntu/c3-pro-idm-frontend/idm
environment = NODE_ENV=prod
command     = /usr/bin/npm 
user        = www-data
autostart   = true
autorestart = true
